# luci-app-simple2fa/root/etc/init.d/simple2fa
#!/bin/sh /etc/rc.common
# Simple2FA - init.d script
# 所有操作都会记录到系统日志，使用 logread | grep simple2fa 查看

START=99
STOP=10

# CGI 路径
CGI_TARGET="/www/cgi-bin/luci"
CGI_SOURCE="/usr/share/luci-app-simple2fa/luci"

# 注入源文件
INJECTION_SOURCE="/usr/share/luci-app-simple2fa/sysauth_inject.html"

log() {
	logger -t simple2fa "[init.d] $1"
	echo "[simple2fa] $1"
}

# 备份并替换文件 (用于 CGI 脚本)
backup_replace() {
	local src="$1" tgt="$2" desc="$3" need_exec="$4"

	if [ ! -f "$src" ]; then
		log "$desc: 跳过 (源文件不存在: $src)"
		return 1
	fi

	if [ ! -e "$tgt" ] && [ ! -e "${tgt}.bak" ]; then
		log "$desc: 跳过 (目标不存在: $tgt)"
		return 1
	fi

	# 备份原始文件（仅首次）
	if [ -e "$tgt" ] && [ ! -e "${tgt}.bak" ]; then
		if ! grep -q "Simple2FA" "$tgt" 2>/dev/null; then
			cp -af "$tgt" "${tgt}.bak"
			log "$desc: 备份成功 -> ${tgt}.bak"
		else
			log "$desc: 跳过备份 (已是我们的文件)"
		fi
	fi

	# 替换文件
	cp -f "$src" "$tgt"
	[ "$need_exec" = "exec" ] && chmod +x "$tgt"
	log "$desc: 替换成功 ($src -> $tgt)"
	return 0
}

# 注入 2FA 脚本到模板文件
inject_template() {
	local tgt="$1" desc="$2"

	if [ ! -f "$INJECTION_SOURCE" ]; then
		log "$desc: 错误 (注入源文件不存在: $INJECTION_SOURCE)"
		return 1
	fi

	if [ ! -e "$tgt" ] && [ ! -e "${tgt}.bak" ]; then
		log "$desc: 跳过 (目标不存在: $tgt)"
		return 1
	fi

	# 备份原始文件（仅首次）
	if [ -e "$tgt" ] && [ ! -e "${tgt}.bak" ]; then
		if ! grep -q "Simple2FA: Runtime injection" "$tgt" 2>/dev/null; then
			cp -af "$tgt" "${tgt}.bak"
			log "$desc: 备份成功 -> ${tgt}.bak"
		else
			log "$desc: 跳过备份 (已包含注入标记)"
		fi
	fi

	# 注入文件：从 .bak 重新生成，保证只注入一次
	if [ -f "${tgt}.bak" ]; then
		cp -af "${tgt}.bak" "$tgt"
		cat "$INJECTION_SOURCE" >> "$tgt"
		log "$desc: 注入成功 ($tgt)"
	else
		log "$desc: 注入失败 (无备份文件)"
	fi
	return 0
}

# 恢复文件
restore_file() {
	local tgt="$1" desc="$2" need_exec="$3"

	if [ -f "${tgt}.bak" ]; then
		rm -f "$tgt"
		cp -af "${tgt}.bak" "$tgt"
		rm -f "${tgt}.bak"
		[ "$need_exec" = "exec" ] && chmod +x "$tgt"
		log "$desc: 恢复成功 (从 .bak)"
	fi
}

start() {
	log "========== START 被调用 =========="

	local enabled=$(uci -q get simple2fa.global.enabled)
	log "配置: enabled='$enabled'"

	local secret=$(uci -q get simple2fa.global.secret)
	if [ -n "$secret" ]; then
		log "密钥: 已配置 (长度=${#secret})"
	else
		log "密钥: 未配置"
	fi

	if [ "$enabled" != "1" ]; then
		log "2FA 未启用 -> 调用 stop 恢复"
		stop
		log "========== START 完成 (未激活) =========="
		return 0
	fi

	# === 1. 替换 CGI 脚本 ===
	backup_replace "$CGI_SOURCE" "$CGI_TARGET" "CGI" "exec"

	# === 2. 注入所有主题的 HTM 模板 ===
	for p in /usr/lib/lua/luci/view/themes/*/sysauth.htm /usr/lib/lua/luci/view/sysauth.htm; do
		if [ -e "$p" ] || [ -e "${p}.bak" ]; then
			inject_template "$p" "HTM模板"
		fi
	done

	# === 3. 注入所有主题 of UT 模板 ===
	for p in /usr/share/ucode/luci/template/themes/*/sysauth.ut \
	         "/usr/lib/ucode/luci/template/themes/*/sysauth.ut" \
	         "/usr/share/ucode/luci/template/sysauth.ut"; do
		if [ -e "$p" ] || [ -e "${p}.bak" ]; then
			inject_template "$p" "UT模板"
		fi
	done

	log "========== START 完成 (2FA 已激活) =========="
}

stop() {
	log "========== STOP 被调用 =========="

	# === 1. 恢复 CGI ===
	restore_file "$CGI_TARGET" "CGI" "exec"

	# === 2. 恢复所有 HTM 模板 ===
	for p in /usr/lib/lua/luci/view/themes/*/sysauth.htm /usr/lib/lua/luci/view/sysauth.htm; do
		restore_file "$p" "HTM模板"
	done

	# === 3. 恢复所有 UT 模板 ===
	for p in /usr/share/ucode/luci/template/themes/*/sysauth.ut \
	         "/usr/lib/ucode/luci/template/themes/*/sysauth.ut" \
	         "/usr/share/ucode/luci/template/sysauth.ut"; do
		restore_file "$p" "UT模板"
	done

	log "========== STOP 完成 =========="
}

restart() {
	log "========== RESTART 被调用 =========="
	stop
	start
}
